# runit/Library/LaunchAgents/local.runit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.runit</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/runsvdir</string>
      <string>-P</string>
      <string>/Users/chberg/.sv</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/chberg/Library/Logs/runit/runsvdir-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chberg/Library/Logs/runit/runsvdir-err.log</string>
  </dict>
</plist>
